# bmc/assemble-lids.service.in
[Unit]
Description=Assemble code update image

[Service]
Type=oneshot
RemainAfterExit=no
ExecStart=/usr/bin/gen-lids-image-tar

