# mmc/xyz.openbmc_project.Software.Version.service.in
[Unit]
Description=Phosphor Version Software Manager

[Service]
ExecStartPre=/bin/sh -c 'mkdir -p /var/tmp/images'
ExecStartPre=/bin/sh -c 'ln -s /var/tmp/images /tmp/images'
ExecStart=/usr/bin/phosphor-version-software-manager
Restart=always
Type=dbus
BusName={BUSNAME}

[Install]
WantedBy=multi-user.target
